# infra/roles/services/letsencrypt/templates/renew-ssl.timer.j2
[Unit]
Description=Timer for renew SSL
Requires=renew-ssl.service

[Timer]
Persistent=true
OnUnitActiveSec=1d
Unit=renew-ssl.service

[Install]
WantedBy=timers.target
